# systemd/container.health.service
[Unit]
Description=conatiner health check job

[Service]
Type=oneshot
ExecStart=/bin/bash <absolute-path-to>/monit_docker_check.sh